# rootfs/etc/services.d/10-php-fpm/run
#!/usr/bin/with-contenv bash

if ${VERBOSE_INIT}
then
    echo "Starting php-fpm..."
fi

exec /usr/sbin/php-fpm -F
ln -s /run/php-fpm${PHP_VERSION}/php-fpm.sock /run/php/php-fpm.sock
