# rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bash

echo "Starting nginx..."
exec /usr/sbin/nginx -c /etc/nginx/nginx.conf
