# rootfs/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv bash

echo "Starting php-fpm..."
exec /usr/sbin/php-fpm -F
ln -s /run/php-fpm${PHP_VERSION}/php-fpm.sock /run/php/php-fpm.sock
